# tools/provisioning/dev/profile.d/profile-path-maven.sh
export PATH=${PATH}:/opt/apache-maven-3.6.3/bin