# tools/provisioning/dev/profile.d/profile-environment-variables.sh
export GEM_HOME=/usr/local/bundle
export AWS_DEFAULT_REGION="eu-west-2"